# crontab
* * * * * /symfony_version.sh > /var/www/public/metrics
